# templates/etc/systemd/system/spigot.socket.j2
[Unit]
Description=Spigot Server Console

[Socket]
ListenFIFO=/run/spigot/console.socket
SocketUser=spigot
SocketGroup=spigot
SocketMode=0600

[Install]
WantedBy=sockets.target
